# supervisor.d/polybar.conf
[program:statusbar]
command=polybar main
autostart=true
