# supervisor.d/urxvtd.conf
[program:urxvtd]
command=urxvtd -q -o
autostart=false
